# react_agent/medical-agent-api.service
[Unit]
Description=Medical Agent Flask API
After=network.target

[Service]
Type=simple
User=participant
WorkingDirectory=__APP_DIR__
EnvironmentFile=__APP_DIR__/.env
ExecStart=__GUNICORN__ \
    --bind 127.0.0.1:3001 \
    --workers 2 \
    --worker-class gthread \
    --threads 4 \
    --timeout 300 \
    --keep-alive 5 \
    server:app
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
